# systemd/audio-api.service
[Unit]
Description=Audio API Server
Documentation=https://github.com/your-org/agent-edge-rs
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
ExecStart=%h/.local/bin/audio_api --bind 0.0.0.0:50051 --max-connections 5
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audio-api

# Security
NoNewPrivileges=true
ProtectHome=false
ProtectSystem=strict
PrivateTmp=true

# Resource limits
MemoryHigh=100M
MemoryMax=200M

# Environment
Environment="RUST_LOG=info"

[Install]
WantedBy=default.target 